# server/aesdsocket.service
[Unit]
Description=aesdsocket

[Service]
ExecStart=/usr/bin/aesdsocket -d

[Install]
WantedBy=multi-user.target

